# systemd/trend-ch4.timer
[Unit]
Description=CH4 GitHub Trending trend report (Mon 08:00)

[Timer]
OnCalendar=Mon *-*-* 08:00:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
