# systemd/trend-ch1.service
[Unit]
Description=Run CH1 Claude / Anthropic trend report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__SYSTEM_DIR__
EnvironmentFile=__HOME__/.env.agent
Environment="PATH=__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin:/snap/bin"
ExecStart=__SYSTEM_DIR__/scripts/run.sh claude-anthropic
StandardOutput=append:__SYSTEM_DIR__/logs/trend-report.log
StandardError=append:__SYSTEM_DIR__/logs/trend-report.log
TimeoutStartSec=1800
