# systemd/trend-ranking.service
[Unit]
Description=Run rolling 4-week Top 5 ranking generation (AI selection)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__SYSTEM_DIR__
EnvironmentFile=__HOME__/.env.agent
Environment="PATH=__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin:/snap/bin"
ExecStart=__SYSTEM_DIR__/scripts/generate-ranking.sh
StandardOutput=append:__SYSTEM_DIR__/logs/ranking.log
StandardError=append:__SYSTEM_DIR__/logs/ranking.log
TimeoutStartSec=900
