# systemd/trend-ch3.service
[Unit]
Description=Run CH3 AI Trends trend report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__SYSTEM_DIR__
EnvironmentFile=__HOME__/.env.agent
ExecStart=__SYSTEM_DIR__/scripts/run.sh ai-trends
StandardOutput=append:__SYSTEM_DIR__/logs/trend-report.log
StandardError=append:__SYSTEM_DIR__/logs/trend-report.log
TimeoutStartSec=1800
